# services/CheckServices.service
[Unit]
Description=Modem Status Checker
After=network.target

[Service]
ExecStart=/app/monitoring/checker/check_services.py
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
