# services/imu.service
[Unit]
Description=MPU-6050 IMU Data Logger
Documentation=IMU sampling service for continuous monitoring
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app/monitoring/imu
ExecStart=/usr/bin/python3 /app/monitoring/imu/imu.py

# Restart configuration - always restart on failure
Restart=always
RestartSec=5

# Resource limits and security
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Give the service time to start
TimeoutStartSec=30
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=imu-logger

# Environment
Environment="PYTHONUNBUFFERED=1"

# Ensure /dev/shm is available
RequiresMountsFor=/dev/shm

# If the service fails to start, keep retrying
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target

